# Lecture_6/date.timer
[Unit]
Description=Run date service every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target